# units/cron-schedule@.target.in
[Unit]
Description={{ package }} %i target
Documentation=man:systemd.cron(7)
Requires=cron-schedule@%i.service
StopWhenUnneeded=yes
