# units/cron-schedule.timer.in
[Unit]
Description={{ package }} {{ schedule }} timer
Documentation=man:systemd.cron(7)
PartOf=cron.target
RefuseManualStart=yes
RefuseManualStop=yes

[Timer]
{{#if persistent}}
Persistent=true
{{/if}}
OnCalendar={{ schedule }}
Unit=cron-{{ schedule }}.target
